# db/crontab
10 * * * * cd /app/db/ && mongo < clear-old-rooms
